# role/app/templates/app.service.j2
[Unit]
Description=Run test app

[Service]
ExecStart=/opt/app.py
User=root
WorkingDirectory=/opt/

[Install]
WantedBy=default.target
